# packaging/core/debian/prerm
#!/bin/bash

/etc/init.d/gridway-core stop >/dev/null 2>&1
update-rc.d -f gridway-core remove
