# packaging/core/debian/prerm
#!/bin/bash

/etc/init.d/gridway-core stop >/dev/null 2>&1
/sbin/chkconfig --del gridway-core
