# packaging/core/debian/prerm
#!/bin/bash

/etc/init.d/gwd stop >/dev/null 2>&1
/sbin/chkconfig --del gwd
