# packaging/core/prerm
#!/bin/bash

/sbin/service gridway-core stop >/dev/null 2>&1
/sbin/chkconfig --del gridway-core
